# supervisord.conf
[supervisord]
nodaemon=true

[program:gogs]
command=su git -c "/home/git/go/src/github.com/gogits/gogs/gogs web"
